# alsactl/alsa-store.service.in
[Unit]
Description=Store Sound Card State
DefaultDependencies=no
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=@sbindir@/alsactl store
StandardOutput=syslog
